# deploy/phansora-api.service.d/cuda-env.conf
# ⚠ LEGACY — NO LONGER REQUIRED by the CosyVoice2 TTS engine.
# This drop-in existed for the old IndexTTS2 + DeepSpeed path, which JIT-compiled fused
# CUDA ops at model-load time and therefore needed nvcc / the CUDA toolkit + a modern gcc.
# CosyVoice2's vLLM and TensorRT wheels ship precompiled CUDA ops, so the API process needs
# only the NVIDIA driver — not the toolkit. Safe to delete this file and skip the copy step;
# kept for reference in case a future engine needs nvcc again.
#
# systemd drop-in for the phansora-api service — makes nvcc / the CUDA toolkit visible
# to the API process so DeepSpeed (INDEXTTS2_USE_DEEPSPEED) and the BigVGAN CUDA kernel
# (INDEXTTS2_USE_CUDA_KERNEL) can JIT-compile their fused ops at model-load time.
#
# The PATH below MUST include the venv's bin (…/.venv/bin) — DeepSpeed's op builder shells
# out to `ninja`, which pip installs INTO the venv, and systemd does not activate the venv
# (so its bin is otherwise absent from PATH). Missing it => "Unable to JIT load the
# transformer_inference op due to ninja not being installed". Adjust the venv path if your
# deployment differs from /var/www/phansora-api/.venv.
#
# Install (after running scripts/install-cuda-toolkit.sh):
#     cp deploy/phansora-api.service.d/cuda-env.conf /etc/systemd/system/phansora-api.service.d/
#     systemctl daemon-reload && systemctl restart phansora-api
#
# RHEL8's default GCC (8.5) is too old to compile torch 2.8's extensions (needs >= 9), so
# builds use gcc-toolset-12 (installed by scripts/install-cuda-toolkit.sh). CC/CXX + the
# toolset bin on PATH make DeepSpeed's builder and nvcc's host compiler use it; the
# toolset lib64 on LD_LIBRARY_PATH provides the matching libstdc++ at .so-load time. On a
# host whose default gcc is already >= 9 (e.g. RHEL9), remove the three gcc-toolset lines.
[Service]
Environment=CUDA_HOME=/usr/local/cuda-12.6
Environment=CC=/opt/rh/gcc-toolset-12/root/usr/bin/gcc
Environment=CXX=/opt/rh/gcc-toolset-12/root/usr/bin/g++
Environment=PATH=/opt/rh/gcc-toolset-12/root/usr/bin:/usr/local/cuda-12.6/bin:/var/www/phansora-api/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=LD_LIBRARY_PATH=/opt/rh/gcc-toolset-12/root/usr/lib64
