# deploy/phansora-api.service.d/cuda-env.conf
# systemd drop-in for the phansora-api service — makes nvcc / the CUDA toolkit visible
# to the API process so DeepSpeed (INDEXTTS2_USE_DEEPSPEED) and the BigVGAN CUDA kernel
# (INDEXTTS2_USE_CUDA_KERNEL) can JIT-compile their fused ops at model-load time.
#
# The PATH below MUST include the venv's bin (…/.venv/bin) — DeepSpeed's op builder shells
# out to `ninja`, which pip installs INTO the venv, and systemd does not activate the venv
# (so its bin is otherwise absent from PATH). Missing it => "Unable to JIT load the
# transformer_inference op due to ninja not being installed". Adjust the venv path if your
# deployment differs from /var/www/phansora-api/.venv.
#
# Install (after running scripts/install-cuda-toolkit.sh):
#     cp deploy/phansora-api.service.d/cuda-env.conf /etc/systemd/system/phansora-api.service.d/
#     systemctl daemon-reload && systemctl restart phansora-api
[Service]
Environment=CUDA_HOME=/usr/local/cuda-12.6
Environment=PATH=/usr/local/cuda-12.6/bin:/var/www/phansora-api/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
